# docker-config/php/crontab/cron
* * * * * /var/shell/notify.sh
* * * * * /var/shell/test.sh
